# devbox/supervisor/conf.d/bootstrap.conf
[program:bootstrap]
user = uesyn
environment = USER="uesyn",HOME="/home/uesyn"
command = /usr/local/share/bootstrap
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
